# Bootloader.PlayStationClassic/payload/recoverykit/systemctl/ftp@.service
[Unit]
Description=FTP Server

[Service]
ExecStart=-/usr/bin/busybox-bs ftpd -w /
StandardInput=socket