# Bootloader.PlayStationClassic/payload/recoverykit/systemctl/telnet@.service
[Unit]
Description=Telnet Server

[Service]
ExecStart=-/usr/bin/busybox-bs telnetd -i
StandardInput=socket