# Bootloader.PlayStationClassic/payload/telnet_ftp_setup/ftp@.service
[Unit]
Description=FTP Server
After=bleemsync.service
Before=sonyapp.service

[Service]
ExecStart=-/usr/bin/busybox-bs ftpd -w /
StandardInput=socket