# Bootloader.PlayStationClassic/payload/telnet_ftp_setup/telnet@.service
[Unit]
Description=Telnet Server
After=bleemsync.service
Before=sonyapp.service

[Service]
ExecStart=-/usr/bin/busybox-bs telnetd -i
StandardInput=socket